# deploy/systemd/transcendence-memory-backend.service
[Unit]
Description=Transcendence Memory Backend (task-rag server)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/transcendence-memory-server
EnvironmentFile=/etc/transcendence-memory/transcendence-memory-backend.env
ExecStart=/usr/bin/env bash ./scripts/run_task_rag_server.sh
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
